# service/99iacgamepad.sh
#!/bin/bash
### BEGIN INIT INFO
# Provides:          insert a coin
# Required-Start:    
# Required-Stop:     
# Default-Start:     
# Default-Stop:      
# Short-Description: iac gamepad
# Description:       Run Insert a Coin gamepad service
### END INIT INFO

# Carry out specific functions when asked to by the system
/bin/iacgamepad

exit 0